# src/rootfs/etc/profile.d/aqua.sh
# aqua: system-wide shell extension

# Code blow is deprecated due to Nix Shell Compliance service.
#export PATH="/opt/aqua/sys/applications:/opt/aqua/applications:/opt/aqua/sys/sbin:/opt/aqua/sbin:$PATH"
