# src/rootfs/etc/profile.d/aqua.sh
# aqua: system-wide shell extension
export PATH="/opt/aqua/sys/applications:/opt/aqua/applications:/opt/aqua/sys/sbin:/opt/aqua/sbin:$PATH"
