# root-skeleton/etc/profile.d/aqua.sh
# aqua: system-wide shell extension

# Code blow is deprecated due to Nix Shell Compliance service.
#export PATH="{{SYS_APPS}}:{{USER_APPS}}:{{SYS_CMDS}}:{{CMDS}}:$PATH"
