# src/features/snapshot/patch/etc/xdg/autostart/me.hysong.aqua.service.SnapshotBootDetector.desktop
[Desktop Entry]
Name=SnapshotBootDetector
GenericName=Snapshot Boot Detector Service
Comment=Detects if current boot session is a snapshot or not.
Exec=/opt/aqua/features/snapshot/snapshot-bootstateutil.sh
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
